# deploy/restore-drill.timer
[Unit]
Description=Weekly zync restore drill

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
